# md2htm.service
[Unit]
Description=Markdown to HTML parser daemon

[Service]
ExecStart=/usr/bin/md2htm daemon

[Install]
WantedBy=multi-user.target
